# thirsty-ai-builder/deploy/ollama.service
# systemd unit for Ollama, hardened for ThirstyAI Builder hosting.
#
# Drop into /etc/systemd/system/ollama.service on the host, then:
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ollama
#   sudo systemctl status ollama
#
# This unit runs Ollama as a dedicated unprivileged user with its own
# data dir under /var/lib/ollama. It binds the API to 127.0.0.1 only
# by default; expose it to the ThirstyAI Builder backend over the
# loopback interface inside the same host, or over a private overlay
# network (Tailscale / WireGuard) when the backend runs elsewhere.

[Unit]
Description=Ollama Service (ThirstyAI Builder)
Documentation=https://ollama.com
After=network-online.target
Wants=network-online.target

[Service]
Type=exec
ExecStart=/usr/local/bin/ollama serve
Restart=always
RestartSec=3
User=ollama
Group=ollama

# Data + model storage. ~5 GB per 7B model, ~30 GB for 70B.
Environment="HOME=/var/lib/ollama"
Environment="OLLAMA_MODELS=/var/lib/ollama/models"
Environment="XDG_CACHE_HOME=/var/lib/ollama/.cache"
Environment="XDG_CONFIG_HOME=/var/lib/ollama/.config"
WorkingDirectory=/var/lib/ollama

# Hardening. Ollama needs a few capabilities for GPU passthrough on
# Linux; everything else is dropped.
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=false
DevicePolicy=closed
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=false
SystemCallArchitectures=native
SystemCallFilter=@system-service
SystemCallErrorNumber=EPERM
UMask=0077

# Resource limits. Adjust MemoryMax to your box. The default 16 GiB
# caps a 7B q4 model; 32 GiB covers a 13B / 33B comfortably.
LimitNOFILE=65536
MemoryMax=32G
TasksMax=4096

# Ollama binds its own port; 11434 is the upstream default.
AmbientCapabilities=CAP_NET_BIND_SERVICE

[Install]
WantedBy=multi-user.target
